# misc/auto_power_profile/99-power-mode.rules
SUBSYSTEM=="power_supply", ATTR{online}=="1", RUN+="/usr/local/bin/power_mode.sh"
SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/usr/local/bin/power_mode.sh"
